# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:plot]
command=plotman plot
redirect_stderr=true

[program:archive]
command=plotman archive
redirect_stderr=true
