# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:plot]
command=plotman plot

[program:archive]
command=plotman archive
